# packaging/worker/systemd/ultimate-terminal-worker.service
[Unit]
Description=Ultimate Terminal Worker
Documentation=https://github.com/stevenvo780/ultimate-terminal
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
EnvironmentFile=/etc/ultimate-terminal/worker.env
ExecStart=/usr/bin/ultimate-terminal-worker
Restart=always
RestartSec=5

# Security hardening (relaxed for terminal access)
NoNewPrivileges=false
ProtectSystem=full
ProtectHome=false
PrivateTmp=false
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true

# Allow PTY allocation
DeviceAllow=/dev/ptmx rw
DeviceAllow=/dev/pts/* rw

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ut-worker

[Install]
WantedBy=multi-user.target
